# .agents/launchd/com.newsicatv.director.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.newsicatv.director</string>
  <key>WorkingDirectory</key>
  <string>/Users/giovannifiore/Desktop/newsica_tv</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/giovannifiore/Desktop/newsica_tv/venv/bin/python</string>
    <string>-u</string>
    <string>/Users/giovannifiore/Desktop/newsica_tv/src/director.py</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/giovannifiore/Desktop/newsica_tv/tmp/director.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/giovannifiore/Desktop/newsica_tv/tmp/director.log</string>
</dict>
</plist>
